# deploy/systemd/odoo-dev.service
[Unit]
Description=Odoo 19 Enterprise - Development (Skymedic)
After=postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=omegaserver02
Group=omegaserver02
WorkingDirectory=/opt/odoo/src/odoo
ExecStart=/opt/odoo/dev/venv/bin/python /opt/odoo/src/odoo/odoo-bin --config=/opt/odoo/dev/config/odoo-dev.conf
Restart=on-failure
RestartSec=5
Environment=PYTHONPATH=/opt/odoo/src/odoo

[Install]
WantedBy=multi-user.target
